# dot_local/misc/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
Environment=XDG_SESSION_TYPE=wayland
ExecStart=
ExecStart=-/sbin/agetty --skip-login --nonewline --noissue --autologin vbt --noclear %I $TERM
